# etc/supervisord.d/hdfs-namenode.ini
[program:hdfs-namenode]
priority=910
command=/opt/qnib/hdfs/namenode/bin/start.sh
redirect_stderr=true
exitcodes=0
stdout_logfile=/var/log/supervisor/%(program_name)s.log
