# etc/supervisord.d/hdfs-configure.ini
[program:hdfs-configure]
priority=999
command=/opt/qnib/hdfs/bin/configure.sh
startsecs=2
exitcodes=0
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
