# etc/supervisord.d/hdfs-datanode.ini
[program:hdfs-datanode]
priority=920
command=/opt/qnib/hdfs/datanode/bin/start.sh
redirect_stderr=true
startsecs=1
exitcodes=0
stdout_logfile=/var/log/supervisor/%(program_name)s.log
