# etc/supervisord.d/hosts-update.ini
[program:hosts-update]
priority=910
command=consul-template -consul localhost:8500 -wait 5s -template "/etc/consul-templates/hosts.ctmpl:/etc/hosts"
redirect_stderr=true
exitcodes=0
stdout_logfile=/var/log/supervisor/%(program_name)s.log
