# etc/supervisord.d/hosts-update.ini
[program:hosts-update]
command=/opt/qnib/hosts-update/bin/start.sh
redirect_stderr=true
exitcodes=0
stdout_logfile=/var/log/supervisor/%(program_name)s.log
